# engine/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:rpcbind]
priority=10
command=/sbin/rpcbind -f -w
autostart=true
autorestart=true

[program:smbd]
priority=20
command=/usr/sbin/smbd -F --no-process-group
autostart=true
autorestart=true

[program:nmbd]
priority=20
command=/usr/sbin/nmbd -F --no-process-group
autostart=true
autorestart=true

[program:vsftpd]
priority=30
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true

[program:nfs]
priority=40
command=/bin/bash -c "rpc.nfsd 8 && exec /usr/sbin/rpc.mountd -F"
autostart=true
autorestart=true

[program:bytebay-engine]
priority=50
command=/usr/local/bin/bytebay-engine -socket /var/bytebay/sockets/engine.sock
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
