# engine/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:dbus]
priority=5
command=/usr/bin/dbus-daemon --system --nofork --nopidfile
autostart=true
autorestart=true

[program:rpcbind]
priority=10
command=/sbin/rpcbind -f -w
autostart=true
autorestart=true

[program:smbd]
priority=20
command=/usr/sbin/smbd -F --no-process-group
autostart=true
autorestart=true

[program:nmbd]
priority=20
command=/usr/sbin/nmbd -F --no-process-group
autostart=true
autorestart=true

[program:vsftpd]
priority=30
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true

[program:vsftpd-ipv6]
priority=30
command=/etc/bytebay/start-vsftpd-ipv6.sh
autostart=true
autorestart=unexpected
startsecs=1
exitcodes=0

[program:ganesha]
priority=40
command=/usr/bin/ganesha.nfsd -F -L /dev/stderr -f /var/lib/bytebay/ganesha.conf
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bytebay-engine]
priority=50
command=/usr/local/bin/bytebay-engine -socket /var/bytebay/sockets/engine.sock
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
